# App-Example/bin/example_cleanup.cron
#!/usr/bin/env perl6

use v6;

# do some cleanup

...;
